# deploy/supervisor.conf
; Supervisor untuk queue worker REKA — §11.5.
; queue: ai (penjanaan draf) + default (mail, whatsapp, dll).
; Salin ke /etc/supervisor/conf.d/reka-worker.conf
;
; PENTING (§Fasa 13): saluran draf HTML memanggil AI 2× + retry — job boleh ambil
; masa lama. Set DB_QUEUE_RETRY_AFTER=360 dalam .env (mesti > --timeout=300 di bawah)
; supaya worker tidak menuntut semula job yang masih berjalan.

[program:reka-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/reka/artisan queue:work --queue=ai,default --tries=1 --timeout=300 --sleep=3 --max-time=3600
directory=/var/www/reka
autostart=true
autorestart=true
stopwaitsecs=310
user=reka
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/reka/storage/logs/worker.log
stopasgroup=true
killasgroup=true
